# ai-ecosystem/config/supervisord.conf
; supervisord.conf — Process Guard (Layer 5)
; Keeps agents alive during session (does NOT survive Codespaces restarts)
; 
; Install: pip install supervisord OR apt-get install supervisor
; Usage:  supervisord -c supervisord.conf
; Check:  supervisorctl -c supervisord.conf status
; Stop:   supervisorctl -c supervisord.conf stop all

[supervisord]
nodaemon=false
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/tmp

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; ── Model Server: Manager/Personality ─────────────────────
[program:llama-manager]
command=bash -c 'export LD_LIBRARY_PATH=$HOME/bin:$LD_LIBRARY_PATH; $HOME/bin/llama-server -m $HOME/ai-models/Qwen3.5-4B-Q4_K_M.gguf -c 8192 --host 127.0.0.1 --port 8080 --threads 2 --no-webui'
autorestart=true
startsecs=10
stopwaitsecs=10
stdout_logfile=/tmp/llama-manager.log
stderr_logfile=/tmp/llama-manager-error.log
priority=10

; ── Model Server: Coder ───────────────────────────────────
[program:llama-coder]
command=bash -c 'export LD_LIBRARY_PATH=$HOME/bin:$LD_LIBRARY_PATH; $HOME/bin/llama-server -m $HOME/ai-models/qwen2.5-coder-7b-instruct-q4_k_m.gguf -c 8192 --host 127.0.0.1 --port 8081 --threads 2 --no-webui'
autorestart=true
startsecs=15
stopwaitsecs=10
stdout_logfile=/tmp/llama-coder.log
stderr_logfile=/tmp/llama-coder-error.log
priority=20

; ── Model Server: Researcher ──────────────────────────────
[program:llama-researcher]
command=bash -c 'export LD_LIBRARY_PATH=$HOME/bin:$LD_LIBRARY_PATH; $HOME/bin/llama-server -m $HOME/ai-models/DeepSeek-R1-Distill-Qwen-7B-Q4_K_M.gguf -c 8192 --host 127.0.0.1 --port 8082 --threads 2 --no-webui'
autorestart=true
startsecs=15
stopwaitsecs=10
stdout_logfile=/tmp/llama-researcher.log
stderr_logfile=/tmp/llama-researcher-error.log
priority=20

; ── Model Server: Memory ───────────────────────────────────
[program:llama-memory]
command=bash -c 'export LD_LIBRARY_PATH=$HOME/bin:$LD_LIBRARY_PATH; $HOME/bin/llama-server -m $HOME/ai-models/Mistral-Nemo-2407-12B-Thinking-Claude-Gemini-GPT5.2-Uncensored-HERETIC.Q4_K_M.gguf -c 8192 --host 127.0.0.1 --port 8083 --threads 2 --no-webui'
autorestart=true
startsecs=20
stopwaitsecs=10
stdout_logfile=/tmp/llama-memory.log
stderr_logfile=/tmp/llama-memory-error.log
priority=20

; ── Keepalive ──────────────────────────────────────────────
[program:keepalive]
command=bash -c 'while true; do date > /tmp/.heartbeat; sleep 300; done'
autorestart=true
startsecs=1
stdout_logfile=/tmp/keepalive.log
priority=5
